# stacks/coding/systemd/coding-tools.target
[Unit]
Description=Coding tools user systemd domain
Wants=codex-cli-update.timer
After=network-online.target

[Install]
WantedBy=default.target
