# stacks/coding/systemd/opencode-update.timer
[Unit]
Description=opencode 6h auto-update
PartOf=coding-tools.target

[Timer]
OnCalendar=*-*-* 00/6:00:00
Persistent=true

[Install]
WantedBy=timers.target
