# stacks/coding/systemd/opencode.service
[Unit]
Description=opencode server (native, home server)
After=network-online.target
Wants=network-online.target
PartOf=coding-tools.target

[Service]
Type=simple
WorkingDirectory=%h
EnvironmentFile=%h/.config/opencode/opencode.env
Environment="PATH=%h/.opencode/bin:%h/.local/bin:%h/.local/nodejs/node-v24.16.0-linux-x64/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin"
ExecStart=%h/.opencode/bin/opencode web --hostname 127.0.0.1 --port 4096
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
