# stacks/codex-lb/systemd/github-oidc-broker-cleanup.timer
[Unit]
Description=Run GitHub OIDC broker API key cleanup every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=5s
Persistent=true

[Install]
WantedBy=timers.target
