# stacks/nvidia-build-lb/systemd/agent-apps-delayed-update.service.d/nblb-cutover-lock.conf
[Service]
ExecStart=
ExecStart=/usr/local/libexec/nvidia-build-lb-agent-apps-delayed-update
