# stacks/forgejo-runner/systemd/forgejo-runner-auth-refresh.service
[Unit]
Description=Refresh shared Forgejo runner Codex ChatGPT auth
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/dongwonttuna/Documents/Programming/home-server-infra/stacks/forgejo-runner
Environment=PATH=/home/dongwonttuna/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/dongwonttuna/Documents/Programming/home-server-infra/stacks/forgejo-runner/scripts/refresh-codex-auth.sh
