# stacks/mcp-suite/systemd/mcp-suite.target
[Unit]
Description=MCP suite user systemd domain
Wants=mcp-suite-update.timer

[Install]
WantedBy=default.target
